# systemd/inky-calendar.service
[Unit]
Description=Inky calendar display update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/inky/inky-calendar
ExecStart=/home/inky/inky-calendar/.venv/bin/python /home/inky/inky-calendar/inky_calendar.py
Environment=INKY_CONFIG=/home/inky/inky-calendar/config.yaml

[Install]
WantedBy=multi-user.target
